# linux/systemd/nm3000.service
[Unit]
Description=NM3000
After=network.target
After=syslog.target

[Install]
WantedBy=multi-user.target
Alias=nm3000

[Service]
